# files/timers/wallpaper.service
[Unit]
Description=wallpaper

[Service]
Type=oneshot
ExecStart=wallpaper.sh
WorkingDirectory=/home/stefan/.bin
  