# files/clamav-weekly.timer

[Unit]
Description=Run clamav scan weekly
  
[Timer]
OnCalendar=Wed 12:00
Persistent=true
  
[Install]
WantedBy=timers.target